# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ============================================
; Scraper - 1분 배치 스케줄러 (cron 대신 loop wrapper 사용)
; ============================================
[program:scraper]
command=/app/docker/scraper_loop.sh
directory=/app
user=untatiz
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=30
stdout_logfile=/app/log/scraper.log
stderr_logfile=/app/log/scraper_error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=PYTHONUNBUFFERED="1",HOME="/home/untatiz"

; ============================================
; Flask Web Server - Gunicorn으로 실행
; ============================================
[program:web]
command=gunicorn --bind 0.0.0.0:5000 --workers 2 --threads 4 --timeout 120 --chdir /app/web untatiz_web:app
directory=/app/web
user=untatiz
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stdout_logfile=/app/log/web.log
stderr_logfile=/app/log/web_error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=PYTHONUNBUFFERED="1",FLASK_ENV="production",HOME="/home/untatiz"

; ============================================
; Process Groups
; ============================================
[group:untatiz]
programs=scraper,web
priority=999
